# modules/cylent/templates/docker-machine-env.plist.erb
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
   <key>Label</key>
   <string>com.barkly.docker-machine-env</string>
   <key>ProgramArguments</key>
   <array>
   <string>/bin/bash</string>
   <string>-c</string>
   <string>"/bin/launchctl setenv DOCKER_IP `/usr/local/bin/docker-machine ip <%= @machine_name %>`"</string>
   </array>
   <key>StartInterval</key>
   <true/>
   <key>KeepAlive</key>
   <false/>
</dict>
</plist>
