# modules/cylent/templates/docker-machine.plist.erb
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
   <key>Label</key>
   <string>com.barkly.docker-machine</string>
   <key>ProgramArguments</key>
   <array>
   <string>/usr/local/bin/docker-machine</string>
   <string>start</string>
   <string><%= @machine_name %></string>
   </array>
   <key>RunAtLoad</key>
   <true/>
   <key>KeepAlive</key>
   <false/>
</dict>
</plist>
